# roles/temperature_to_es/templates/usr/lib/systemd/system/temperature_to_es.service.j2
[Unit]
Description=temperature_to_es service

[Service]
Type=oneshot
ExecStart=/usr/bin/python {{ inst_home }}/temperature_to_es.py

[Install]
WantedBy=multi-user.target
