# roles/temperature_to_es/templates/usr/lib/systemd/system/temperature_to_es.timer.j2
[Unit]
Description=temperature_to_es timer

[Timer]
OnCalendar=*:0/{{ frequency_mins }}
AccuracySec=100ms

[Install]
WantedBy=timers.target
